# spx_bot.service
[Unit]
Description=SPX Tracking Telegram Bot
After=network.target

[Service]
# Thay "pi" bằng username bạn dùng để chạy bot
User=pi
Group=pi

# Thay "/home/pi/spx_tracking" bằng đường dẫn tuyệt đối đến thư mục dự án
WorkingDirectory=/home/pi/spx_tracking
ExecStart=/usr/bin/python3 main.py

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target